# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[inet_http_server]
port=127.0.0.1:9001
username=
password=

[supervisorctl]
serverurl=http://127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:ingenium-api]
user=appuser
command=node dist/scripts/api-server.js
directory=/app/services/ingenium-api
autorestart=true
stdout_logfile=/app/.ingenium/logs/%(program_name)s-stdout.log
stdout_logfile_maxbytes=5242880
stdout_logfile_backups=2
stderr_logfile=/app/.ingenium/logs/%(program_name)s-stderr.log
stderr_logfile_maxbytes=5242880
stderr_logfile_backups=2
environment=INGENIUM_HOME="/app/.ingenium",INGENIUM_CORE_DB_PATH="/app/.ingenium/data",CORS_ORIGIN="*"

[program:ingenium-dashboard]
user=appuser
command=node server.js
directory=/app/services/ingenium-dashboard
autorestart=true
stdout_logfile=/app/.ingenium/logs/%(program_name)s-stdout.log
stdout_logfile_maxbytes=5242880
stdout_logfile_backups=2
stderr_logfile=/app/.ingenium/logs/%(program_name)s-stderr.log
stderr_logfile_maxbytes=5242880
stderr_logfile_backups=2
environment=NEXT_PUBLIC_API_URL="http://localhost:4097/api/v1",OPENCODE_SERVER_PASSWORD="%(ENV_OPENCODE_SERVER_PASSWORD)s"

[program:opencode-server]
user=appuser
command=opencode web --port 4096 --hostname 0.0.0.0
directory=/app
autorestart=true
stdout_logfile=/app/.ingenium/logs/%(program_name)s-stdout.log
stdout_logfile_maxbytes=5242880
stdout_logfile_backups=2
stderr_logfile=/app/.ingenium/logs/%(program_name)s-stderr.log
stderr_logfile_maxbytes=5242880
stderr_logfile_backups=2
environment=HOME="/home/appuser"

[program:opencode-iframe]
user=appuser
command=opencode web /workspace --port 4098 --hostname 0.0.0.0
directory=/workspace
autorestart=true
stdout_logfile=/app/.ingenium/logs/%(program_name)s-stdout.log
stdout_logfile_maxbytes=5242880
stdout_logfile_backups=2
stderr_logfile=/app/.ingenium/logs/%(program_name)s-stderr.log
stderr_logfile_maxbytes=5242880
stderr_logfile_backups=2
environment=OPENCODE_SERVER_PASSWORD="",HOME="/home/appuser"
